# content-agent/scheduler/com.simplerabbit.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.simplerabbit.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/leannfrank/Desktop/simple rabbit/content-agent/agent.py</string>
    </array>

    <!-- Runs Monday–Friday at 5:00 AM (uses system timezone — make sure your Mac is set to Eastern) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leannfrank/Desktop/simple rabbit/content-agent</string>

    <key>StandardOutPath</key>
    <string>/Users/leannfrank/Desktop/simple rabbit/content-agent/logs/sage.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/leannfrank/Desktop/simple rabbit/content-agent/logs/sage.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
