# content-agent/scheduler/com.simplerabbit.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    macOS launchd plist — Simple Rabbit Content Agent
    Runs agent.py daily at 7:00 AM.

    INSTALL:
        cp com.simplerabbit.agent.plist ~/Library/LaunchAgents/
        launchctl load ~/Library/LaunchAgents/com.simplerabbit.agent.plist

    UNINSTALL:
        launchctl unload ~/Library/LaunchAgents/com.simplerabbit.agent.plist
        rm ~/Library/LaunchAgents/com.simplerabbit.agent.plist

    VIEW LOGS:
        tail -f /tmp/simplerabbit-agent.log
        tail -f /tmp/simplerabbit-agent-error.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.simplerabbit.agent</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use the full path to your Python 3 interpreter -->
        <!-- Run: which python3  to find the right path -->
        <string>/usr/bin/python3</string>
        <string>/Users/leannfrank/Desktop/simple rabbit/content-agent/agent.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leannfrank/Desktop/simple rabbit/content-agent</string>

    <!-- Run daily at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/simplerabbit-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/simplerabbit-agent-error.log</string>

    <!-- Run even if the scheduled time was missed (e.g. Mac was asleep) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
